# roundcube/rootfs/etc/s6-overlay/s6-rc.d/init-php-fpm/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community App: Roundcube
# Configures PHP-FPM for use with Roundcube
# ==============================================================================

# Generate Ingress configuration
bashio::var.json \
    name "ingress" \
    port "^9002" \
    base "$(bashio::app.ingress_entry)" \
    | tempio \
        -template /etc/php84/templates/php-fpm.gtpl \
        -out /etc/php84/php-fpm.d/ingress.conf

# Generate direct access configuration, if enabled.
if bashio::var.has_value "$(bashio::app.port 80)"; then
    bashio::var.json \
        name "www" \
        port "^9001" \
        | tempio \
            -template /etc/php84/templates/php-fpm.gtpl \
            -out /etc/php84/php-fpm.d/www.conf
fi