# roundcube/rootfs/etc/s6-overlay/s6-rc.d/init-roundcube/run
#!/command/with-contenv bashio
# shellcheck disable=SC2086,SC2016
# ==============================================================================
# Home Assistant App: Roundcube
# Configures Roundcube
# ==============================================================================
export host
export password
export port
export username
export database

host=$(bashio::services "mysql" "host")
password=$(bashio::services "mysql" "password")
port=$(bashio::services "mysql" "port")
username=$(bashio::services "mysql" "username")

# Modify config files
sed -i "s#\['db_dsnw']\ =#['db_dsnw']\ = 'mysql://${username}:${password}@${host}:${port}/roundcubemail';#g" /var/www/roundcube/config/config.inc.php

# Modify config files for S6-logging
#sed -i 's#^ + .*$# + -^auth\\. -^authpriv\\. -mail\\. $T ${dir}/everything#' /etc/s6-overlay/s6-rc.d/syslogd-log/run
#sed -i 's#^ + .*$# + -^auth\\. -^authpriv\\. -mail\\. $T ${dir}/everything#' /run/service/syslogd-log/run.user
#sed -i 's#^backtick .*$#backtick -D "n20 s1000000 T 1" line { printcontenv S6_LOGGING_SCRIPT }#' /etc/s6-overlay/s6-rc.d/syslogd-log/run
#sed -i 's#^backtick .*$#backtick -D "n20 s1000000 T 1" line { printcontenv S6_LOGGING_SCRIPT }#' /run/service/syslogd-log/run.user
#sed -i 's#^s6-socklog .*$#s6-socklog -d3 -U -t3000 -x /run/systemd/journal/dev-log#' /etc/s6-overlay/s6-rc.d/syslogd/run
#sed -i 's#^s6-socklog .*$#s6-socklog -d3 -U -t3000 -x /run/systemd/journal/dev-log#' /run/service/syslogd/run.user

if ! bashio::services.available 'mysql'; then
    bashio::log.fatal \
    "Local database access should be provided by the MariaDB app"
    bashio::exit.nok \
    "Please ensure it is installed and started"
fi

database=$(\
    mariadb \
        -u "${username}" -p"${password}" \
        -h "${host}" -P "${port}" \
        --skip-column-names \
        --skip-ssl \
        -e "SHOW DATABASES LIKE 'roundcubemail';"
)

if ! bashio::var.has_value "${database}"; then
    bashio::log.info "Creating database for Roundcube"
    mariadb \
        -u "${username}" -p"${password}" \
        --skip-ssl \
        -h "${host}" -P "${port}" \
            < /etc/roundcube/createdb.sql
    bashio::log.info "Creating db tables for Roundcube"
    mariadb \
        roundcubemail \
        -u "${username}" -p"${password}" \
        --skip-ssl \
        -h "${host}" -P "${port}" \
            < /var/www/roundcube/SQL/mysql.initial.sql
fi

if bashio::var.has_value "${database}"; then
    bashio::log.info "Making sure Database schema is up to date"
    sed -i "1 s|$|84|" /var/www/roundcube/bin/updatedb.sh
    /var/www/roundcube/bin/updatedb.sh --package=roundcube --dir=/var/www/roundcube/SQL
fi
